# deploy/remotevoice.service
# systemd 服务：RemoteVoice 中继
# 安装：
#   sudo cp deploy/remotevoice.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now remotevoice
#   sudo journalctl -u remotevoice -f   # 看日志
[Unit]
Description=RemoteVoice Relay (pairing + Doubao ASR + LLM)
After=network.target

[Service]
Type=simple
User=remotevoice
Group=remotevoice
WorkingDirectory=/opt/remotevoice/relay
EnvironmentFile=/opt/remotevoice/relay/.env
ExecStart=/usr/bin/node /opt/remotevoice/relay/dist/server.js
Restart=on-failure
RestartSec=3
# 资源限制（可选）
MemoryMax=512M

[Install]
WantedBy=multi-user.target
